# d-i/source/hw-detect/net-hotplug.rules
SUBSYSTEM=="net",	RUN+="/etc/hotplug.d/net/hw-detect.hotplug"
